# dnsmasq/PLIST
@comment $NetBSD$
man/man8/dnsmasq.8
sbin/dnsmasq
share/examples/rc.d/dnsmasq
