# aolserver/PLIST
@comment $NetBSD$
share/examples/rc.d/nsd
