# darkcoin/PLIST
@comment $NetBSD$
bin/darkcoin-qt
sbin/darkcoind
share/examples/darkcoin/darkcoin.conf
share/examples/rc.d/darkcoind
