# bitlbee/PLIST
@comment $NetBSD$
etc/rc.d/bitlbee
man/man8/bitlbee.8
sbin/bitlbee
share/doc/bitlbee/FAQ
share/doc/bitlbee/README
share/doc/bitlbee/TODO
share/bitlbee/data/motd.txt
share/bitlbee/data/help.txt
@dirrm share/bitlbee/data
@dirrm share/doc/bitlbee
