# umurmur/PLIST
@comment $NetBSD$
sbin/umurmurd
share/examples/rc.d/umurmur
share/examples/umurmur/umurmur.conf
