# miniupnpd/PLIST
@comment $NetBSD$
man/man1/miniupnpd.1
sbin/miniupnpd
share/examples/miniupnpd/miniupnpd.conf
share/examples/rc.d/miniupnpd
@dirrm share/examples/miniupnpd
