# cntlm/PLIST
@comment $NetBSD$
bin/cntlm
man/man1/cntlm.1
share/examples/cntlm/cntlm.conf
share/examples/rc.d/cntlm
