# samhain/PLIST
@comment $NetBSD$
sbin/samhain
share/examples/rc.d/samhain
etc/samhainrc.sample
@dirrm var/samhain
