# boinc-client/PLIST
@comment $NetBSD$
bin/boinc_client
bin/boinc_gui
bin/switcher
share/examples/rc.d/boinc
