# boinc-client/PLIST
@comment $NetBSD$
bin/boinc_client
bin/boinc_gui
lib/libbenchmark.a
share/examples/rc.d/boinc
