# tetrinetx/PLIST
@comment $NetBSD$
bin/tetrix
share/examples/rc.d/tetrinetx
share/tetrinetx/game.conf
share/tetrinetx/game.motd
share/tetrinetx/game.pmotd
share/tetrinetx/game.secure
@dirrm share/tetrinetx
