# deforaos-init/PLIST
@comment $NetBSD$
bin/Init
etc/AppInterface/Init.interface
share/examples/rc.d/Init
