# mu-conference/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:33:58 reed Exp $
bin/mu-conference
share/doc/muc/ChangeLog
share/doc/muc/FAQ
share/doc/muc/README
share/doc/muc/TODO
share/doc/muc/XEP0045_SUPPORT
share/examples/jabber/muc-default.xml
share/examples/jabber/style.css
share/examples/rc.d/muc
