# linuxigd-current/PLIST
@comment $NetBSD$
man/man8/upnpd.8
sbin/upnpd
share/examples/linuxigd/dummy.xml
share/examples/linuxigd/gateconnSCPD.xml
share/examples/linuxigd/gatedesc.xml
share/examples/linuxigd/gateicfgSCPD.xml
share/examples/linuxigd/ligd.gif
share/examples/linuxigd/upnpd.conf
share/examples/rc.d/linuxigd
@dirrm share/examples/linuxigd
