# skkserv/PLIST
@comment $NetBSD$
libexec/skkserv
share/examples/rc.d/skkserv
share/skk/README.skkserv
