# cvsd/PLIST
@comment $NetBSD$
etc/rc.d/cvsd
man/man5/cvsd.conf.5
man/man8/cvsd-buildroot.8
man/man8/cvsd-passwd.8
man/man8/cvsd.8
sbin/cvsd
sbin/cvsd-buildroot
sbin/cvsd-passwd
