# dhcpcd/PLIST
@comment $NetBSD$
sbin/dhcpcd
man/man8/dhcpcd.8
share/examples/rc.d/dhcpcd
