# honeyd-arpd/PLIST
@comment $NetBSD$
etc/rc.d/honeydarpd
man/man8/honeyd-arpd.8
sbin/honeyd-arpd
