# honeyd-arpd/PLIST
@comment $NetBSD$
etc/rc.d/arpd
man/man8/arpd.8
sbin/arpd
