# bitcoin/PLIST
@comment $NetBSD$
bin/bitcoin-qt
bin/bitcoind
@comment share/examples/bitcoin/bitcoin.conf
@comment share/examples/rc.d/bitcoind
