# qico/PLIST
@comment $NetBSD$
bin/qico
${QCC}bin/qcc
share/examples/qico/qico.conf
share/examples/rc.d/qico
@dirrm share/examples/qico
