# oak/PLIST
@comment $NetBSD$
bin/oak
man/man1/oak.1
share/examples/rc.d/oak
share/examples/oak/oak.conf
@dirrm share/examples/oak
